# deploy/systemd/silhouette-core.service
[Unit]
Description=Silhouette Agency OS - Core Cognitive Daemon
After=network.target docker.service
Requires=network.target

[Service]
Type=simple
# Update USER and GROUP to the correct deployment user
User=silhouette
Group=silhouette

# Update WorkingDirectory to absolute path of your clone
WorkingDirectory=/opt/Silhouette-Agency-OS-OpenSource

# Environment Variables
Environment=NODE_ENV=production
Environment=PORT=3005

# Command to start the core orchestrator and nervous system
ExecStart=/usr/bin/npm run server

# Restart rules: Always restart if crashed, wait 5 seconds between restarts
Restart=always
RestartSec=5

# Logging: Standard output and error will go to journald
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=silhouette-core

[Install]
WantedBy=multi-user.target
